# supervisor/pilot-web.conf
[program:pilot-web]
command=/bin/bash -c "/opt/pilot-web/Ascon.Pilot.Web"
directory=/opt/pilot-web
user=root
