# supervisor/pilot-server.conf
[program:pilot-server] 
command=/opt/pilot-server/Ascon.Pilot.Daemon /opt/pilot-server/settings.xml
user=root
